# supervisord.conf
[unix_http_server]
file=/tmp/supervisor.sock   ; (the path to the socket file)


[inet_http_server]         ; inet (TCP) server disabled by default
port=*:9001        ; (ip_address:port specifier, *:port for all iface)

[supervisord]
logfile=/tmp/supervisord.log ; (main log file;default $CWD/supervisord.log)
logfile_maxbytes=50MB        ; (max main logfile bytes b4 rotation;default 50MB)
logfile_backups=10           ; (num of main logfile rotation backups;default 10)
loglevel=info                ; (log level;default info; others: debug,warn,trace)
pidfile=/tmp/supervisord.pid ; (supervisord pidfile;default supervisord.pid)
nodaemon=false               ; (start in foreground if true;default false)
minfds=1024                  ; (min. avail startup file descriptors;default 1024)
minprocs=200                 ; (min. avail process descriptors;default 200)

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock ; use a unix:// URL  for a unix socket

[program:resque_subscribers]
command=/usr/bin/php resque
process_name=%(program_name)s_%(process_num)s
numprocs=4
directory=worker
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=logs/resque.log
environment=QUEUE="subscribers,push",REDIS_BACKEND="localhost:6379",PIDFILE=WORKER_SUB_%(process_num)s.PID,WORKER=SUB_%(process_num)s,DELAY=%(process_num)s,INTERVAL=1,APP_INCLUDE="logger.php"

[program:resque_analytics]
command=/usr/bin/php resque
process_name=%(program_name)s_%(process_num)s
numprocs=2
directory=worker
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=logs/resque.log
environment=QUEUE="analytics",REDIS_BACKEND="localhost:6379",PIDFILE=WORKER_ANALYTICS_%(process_num)s.PID,WORKER=A_%(process_num)s,DELAY=%(process_num)s,INTERVAL=10,APP_INCLUDE="logger.php"
